# power-rules-daemon.service
[Unit]
Description=Power Rules Daemon
After=graphical-session.target

[Service]
ExecStart=/usr/bin/power-rules-daemon
Restart=on-failure

[Install]
WantedBy=default.target
